# packaging/linux/deb/prerm
#!/bin/sh

if [ "$1" = "failed-upgrade" ]; then
  exit 1
fi

if [ "$1" = "remove" ]; then
  /etc/init.d/nomnomlog stop
  update-rc.d -f nomnomlog remove >/dev/null
fi

